# client/tests/image/etc/supervisor/conf.d/daemons.conf
[program:nginx]
command = nginx -g 'daemon off;'
autorestart = true
redirect_stderr = true

[program:x11vnc]
command = x11vnc -shared -forever
autorestart = true
redirect_stderr = true

[program:xfce4]
environment = XFCE_PANEL_MIGRATE_DEFAULT=1
command = xvfb-run --auth-file=/root/.Xauthority --server-args='-screen 0 1024x768x24' --server-num 0 startxfce4
autorestart = true
redirect_stderr = true
stopasgroup = true

[program:tomato-server]
environment = TOAMTO_DISABLE_DEBUG_TOOLBAR=1
directory = /app/server
command = /venv/server/bin/python manage.py runserver
autorestart = true
redirect_stderr = true
