# client/tests/image/etc/supervisor/conf.d/xfce4.conf
[program:xfce4]
environment = XFCE_PANEL_MIGRATE_DEFAULT=1
command = xvfb-run --auth-file=/root/.Xauthority --server-args='-screen 0 1024x768x24' --server-num 0 startxfce4
autorestart = true
redirect_stderr = true
